# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=/bin/bash /start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log
